# docksound.rules
ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="0487", ACTION=="add", RUN+="/root/docksound.sh"
